# gateway/systemd/cortex-gateway.service
[Unit]
Description=Cortex Gateway - REST (phone app) + Streamable-HTTP MCP (AI connectors)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/turfptax/cortex-gateway
EnvironmentFile=/home/turfptax/cortex-gateway/.env
ExecStart=/usr/bin/python3 -m cortex_gateway
Restart=on-failure
RestartSec=3
# Pi Zero / small-RAM mitigation parity with cortex-core.
OOMScoreAdjust=-400

[Install]
WantedBy=multi-user.target
